# crontab
0 0 * * 0 python /app/retrain.py >> /var/log/cron.log 2>&1
